# infrastructure/systemd/kelvin-api.service
[Unit]
Description=Kelvin Assistant FastAPI service
Documentation=https://github.com/ZoltanKarika/Kelvin-Assistant
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=kelvin
Group=kelvin
WorkingDirectory=/opt/kelvin-assistant
EnvironmentFile=/etc/kelvin-assistant/kelvin.env
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=KELVIN_SETTINGS_ENV_FILE=/etc/kelvin-assistant/kelvin.env
ExecStart=/opt/kelvin-assistant/.venv/bin/kelvin-api
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGINT

# Restrict the service without blocking its network access or state directory.
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/etc/kelvin-assistant/kelvin.env
StateDirectory=kelvin-assistant
StateDirectoryMode=0750
UMask=0027

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
